# core/profile.d/proxy.sh
export http_proxy=http://http-proxy.kiss.svc.ops.netai-cloud:1080
export https_proxy=$http_proxy
export no_proxy=0,1,2,3,4,5,6,7,8,9,.netai-cloud,localhost,localdomain
