# core/profile.d/proxy.sh
export HTTP_PROXY=http://http-proxy.kiss.svc.ops.netai-cloud:1080
export HTTPS_PROXY=$HTTP_PROXY
export NO_PROXY=0,1,2,3,4,5,6,7,8,9,.netai-cloud,localhost,localdomain
